# scripts/com.codex.watch-notify.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.codex.watch-notify</string>
  <key>ProgramArguments</key><array><string>__PROGRAM__</string></array>
  <key>EnvironmentVariables</key><dict>
    <key>CODEX_WATCH_NTFY_URL</key><string>__URL__</string>
    <key>CODEX_WATCH_NTFY_TOKEN</key><string>__TOKEN__</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>StandardErrorPath</key><string>__LOG__</string>
</dict></plist>
